# rebuild-image.cron
#!/bin/bash

/home/bjg/git/docker-email/build > last-build.log 2>&1
