# spam-learn.cron
#!/bin/bash

/home/bjg/git/docker-email/demail-shell /sbin/spam-learn > /var/log/mail/spam-learn.log 2>&1
